# ufw-rules/cloudflare-ufw.timer
[Unit]
Description=Run cloudflare-ufw every week

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
