# dynamic-dns/cloudflare-dyndns.service
[Unit]
Description=Update Cloudflare DNS records
Documentation=https://github.com/cloudflare-dyndns

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cloudflare-dyndns
